# examples/docker/renewer/crontab
# min	hour	day	month	weekday	command
*	*	*	*	*	step ca renew --force /var/local/step/site.crt /var/local/step/site.key
